# twamp-server/twampd.service
[Unit]
Description=TWAMP Daemon - Two-Way Active Measurement Protocol Server
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/twampd
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10

# Simple logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twampd

[Install]
WantedBy=multi-user.target
